# supervisord.conf
[supervisord]

nodaemon=true

[program:web-app]

command=python manage.py

[program:consumer]

command=python consumer.py
